# packages/azureml-inference-server-http/azureml_inference_server_http-0.7.2-py3-none-any.whl/azureml_inference_server_http/prepost_server/runit/prepost/finish
#!/bin/bash

exit_code="$1" # The exit code from prepost
signal="$2"    # The signal which caused prepost to exit (or 0)

echo "`date -uIns` - prepost/finish $@"
echo "`date -uIns` - Exit code $exit_code is not normal. Killing image."

killall -SIGHUP runsvdir
